# systemd/x3ui-tunnel-2.service
[Unit]
Description=SSH Tunnel to 3x-UI Server 2 (147.45.50.129)
After=network.target

[Service]
Type=simple
User=root
# Туннель: хост:38869 -> VPN сервер:45087 (3x-UI слушает на 127.0.0.1:45087)
ExecStart=/usr/bin/ssh -N -L 0.0.0.0:38869:127.0.0.1:45087 -i /root/fiorevpn/ssh/x3ui_key_2 -o StrictHostKeyChecking=no -o ServerAliveInterval=30 -o ServerAliveCountMax=3 root@147.45.50.129
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

